# server/config/sample.upstart.conf
description "heartBeat Server"
author  "lucas.mouilleron@me.com"

start on runlevel [234]
stop on runlevel [0156]

chdir /root/heartBeat
exec python3 /root/heartBeat/server.py
respawn